# route53-updater.service
[Unit]
Description=Route53 IP Updater
After=network.target

[Service]
Type=oneshot
User=your_username
ExecStart=/path/to/route53-updater.sh
Environment=HOME=/home/your_username

[Install]
WantedBy=multi-user.target